# src/scripts/etc/inittab
# Begin /etc/inittab
::sysinit:/etc/init.d/rcS
tty1::askfirst:/sbin/getty 38400 tty1
::respawn:/sbin/getty -L ttyAMA0 115200 vt100
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
# End /etc/inittab
